# launchd/com.user.claude-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.__USERNAME__.claude-proxy</string>
	<key>ProgramArguments</key>
	<array>
		<string>__HOME__/.claude-code-proxy/.venv/bin/python</string>
		<string>-m</string>
		<string>uvicorn</string>
		<string>app.main:app</string>
		<string>--host</string>
		<string>127.0.0.1</string>
		<string>--port</string>
		<string>3456</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__HOME__/.claude-code-proxy</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/claude-proxy.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-proxy.err</string>
</dict>
</plist>
